# src/conf/supervisor_eos.conf
[program:eos]
command=/start_eos.sh
redirect_stderr=true
