# src/conf/supervisord_eoskeeper.conf
[program:eoskeeper]
command=/start_eoskeeper.sh
redirect_stderr=true
